# himawari.service
[Unit]
Description="Update wallpaper"

[Service]
ExecStart=/home/lia/dotfiles/himawari-wallpaper/invoke.sh
RuntimeMaxSec=90

